# harness/residency/gable-summon.service
# gable-summon.service — systemd USER unit for the WP-H9 summon adapter.
# Runs under res-gable's user manager (lingering already enabled by 02-podman).
# The daemon is HOST-side as res-gable: its run-resident.sh invocations carry
# the res-gable uid the broker (SO_PEERCRED) and nftables wall key on. Its
# egress rides the same res-gable wall: loopback->Disjorn + api.anthropic.com.
#
# Install: harness/keyboard/05-gable.sh (builds the venv this points at).
# Stopping this unit is a plink-side kill switch for summons; so is
# summon.toml (plink-owned config dir).

[Unit]
Description=Gable summon adapter (Disjorn residency, %u)
Wants=network-online.target
After=network-online.target

[Service]
# run-resident.sh env for the per-summon container (same deviations as
# Claudette's unit, same reasons): /home/plink is 0700 so the documented
# defaults are unreachable from a res-* uid — point at the /srv config dir
# and the root-owned harness copies; mount the clean repo mirror read-only.
Environment=RESIDENT_CONFIG_DIR=/srv/disjorn-resident-config/gable
Environment=RESIDENT_HOUSE_MEMORY=/usr/local/lib/disjorn/house_memory
Environment=RESIDENT_DISJORN_RO=/srv/disjorn-ro
ExecStart=/usr/local/lib/disjorn/residency-venv/bin/python \
    /usr/local/lib/disjorn/residency/run_summon.py \
    --config /srv/disjorn-resident-config/res-gable/summon.toml
# The daemon logs to stderr; append it to the file the broker's read-own-log
# serves for res-gable ([residents.res-gable].log_path in broker.toml) —
# without this the configured log path never receives a byte.
StandardOutput=append:/home/res-gable/logs/gable.log
StandardError=append:/home/res-gable/logs/gable.log
Restart=on-failure
RestartSec=15
MemoryMax=1G
TasksMax=128

[Install]
WantedBy=default.target
